# systemd/rss.timer
[Unit]
Description=Fetch rss timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=10m

[Install]
WantedBy=default.target
